# usr/lib/systemd/user/xdg_menu.service
[Unit]
Description=Atualiza xdg_menu
​
[Service]
Type=simple
ExecStart=/usr/bin/xdg_menu
